# install/ansible/roles/app/files/ntag-pi.service
[Unit]
Description=NTAG-Pi

[Service]
ExecStart=/opt/ntag-pi/venv/bin/python /opt/ntag-pi/src/main.py

[Install]
WantedBy=multi-user.target
